# ds360go.service
[Unit]
Description=ds360go service

[Service]
Type=simple
ExecStart=

[Install]
WantedBy=default.target